# etc/init/vagrant-mounted.conf
# start services on vagrant mounted
start on vagrant-mounted
exec sudo service httpd start